# src/wifi/WiFi-Monitor.service
[Unit]
Description=Broadcom T2 WiFi Guardian (Auto-Recovery)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/sbin/WiFi-Monitor daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
